# crontab.cron
0 6,12,18 * * * /home/pi/work/thisnote/go.sh
