# crontab.cron
0 6,12,18,23 * * * /home/pi/work/thisnote/main.py
