# crontab.cron
0 6,12,18,21 * * * /home/pi/work/thisnote/go.sh
